# plugins/session-archiver/templates/launchd-drain.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  macOS spool-drain agent for session-archiver (sync_mode: spool).
  Install:
    1. Edit ProgramArguments path below to your installed plugin location, e.g.
       ~/.claude/plugins/cache/<marketplace>/session-archiver/scripts/drain.sh
    2. Set your own reverse-DNS prefix in the Label below (e.g. com.yourname...) before copying.
    3. cp this file to ~/Library/LaunchAgents/REPLACE_WITH_REVERSE_DNS.session-archiver-drain.plist
    4. launchctl load ~/Library/LaunchAgents/REPLACE_WITH_REVERSE_DNS.session-archiver-drain.plist
  Runs drain.sh every 300s. Logs to <data-dir>/archive.log (the script's own log).

  NOTE: launchd starts agents with a minimal PATH (/usr/bin:/bin:/usr/sbin:/sbin)
  that does NOT include Homebrew, so drain.sh can't find `jq` and aborts. The
  EnvironmentVariables/PATH below prepends the Homebrew bin dirs — /opt/homebrew/bin
  (Apple Silicon) and /usr/local/bin (Intel). Adjust if jq lives elsewhere.
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>REPLACE_WITH_REVERSE_DNS.session-archiver-drain</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>REPLACE_WITH_PATH_TO/session-archiver/scripts/drain.sh</string>
  </array>
  <key>StartInterval</key>
  <integer>300</integer>
  <key>RunAtLoad</key>
  <true/>
  <key>ProcessType</key>
  <string>Background</string>
</dict>
</plist>
